# lapahub/rootfs/etc/s6-overlay/s6-rc.d/lapahub/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# LapaHub Addon service
# ==============================================================================

bashio::log.info "Starting LapaHub Bridge service..."

# Export log level from options
LOG_LEVEL=$(bashio::config 'log_level' 'info')
export LOG_LEVEL

exec python3 /app/run.py
